# telegram-bot.service
[Unit]
Description=Telegram Alarm Bot
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/telegram-bot
Environment="PATH=/home/%i/telegram-bot/venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/%i/telegram-bot/venv/bin/python /home/%i/telegram-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

